# ops/backup/pg-dump.service
[Unit]
Description=Postgres weekly backup
Documentation=file:///home/rod/stack/ops/backup/pg-dump.sh
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
User=rod
ExecStart=/home/rod/stack/ops/backup/pg-dump.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pg-dump
